# cgcs-patch/bin/sw-patch-controller.service
[Unit]
Description=TIS Patching Controller
After=syslog.service network-online.target sw-patch.service
Before=sw-patch-agent.service sw-patch-controller-daemon.service

[Service]
Type=oneshot
User=root
ExecStart=/etc/init.d/sw-patch-controller start
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target

